# infra/ncp/prod/qbook-app.service
[Unit]
Description=Questbook app-api container
After=docker.service network-online.target qbook-bootstrap.service
Requires=docker.service qbook-bootstrap.service

[Service]
EnvironmentFile=/etc/questbook/deploy.env
ExecStartPre=-/usr/bin/docker rm -f qbook-app
ExecStartPre=/usr/bin/docker pull ${REGISTRY}/qbook-app:${IMAGE_TAG}
ExecStart=/usr/bin/docker run --name qbook-app --rm -p 8100:8100 \
  --env-file /etc/questbook/app.env ${REGISTRY}/qbook-app:${IMAGE_TAG}
ExecStop=/usr/bin/docker stop qbook-app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
